# rootfs/etc/s6-overlay/s6-rc.d/rclone-mount-check/up
/etc/s6-overlay/scripts/rclone-mount-check-up
